# tox-easy-bootstrap/debian/cron.d
12	3	*	*	4	root	/usr/bin/tox-easy-bootstrap > /dev/null
